# config/dabradio.service
[Unit]
Description=DAB+ Bluetooth Radio Web Interface
After=network.target bluetooth.target bluealsa.service
Wants=bluealsa.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/dab-radio
ExecStart=/opt/dab-radio/venv/bin/python /opt/dab-radio/server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
